# crontab
*/5  *  *  *  * php -f /var/www/nextcloud/cron.php
*/5  *  *  *  * php /var/www/nextcloud/occ preview:pre-generate
